# scripts/openclaw-token-aggregator.service
[Unit]
Description=OpenClaw Token Aggregator Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /home/wzgrx/.openclaw/channels/feishu/token-aggregator-daemon.js
Restart=on-failure
RestartSec=10
WorkingDirectory=/home/wzgrx/.openclaw/channels/feishu
Environment=NODE_DISABLE_COMPILE_CACHE=1

[Install]
WantedBy=default.target
