# scripts/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (v2026.5.7)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /home/wzgrx/.nvm/versions/node/v24.15.0/lib/node_modules/openclaw/dist/index.js gateway --port 18789
Restart=on-failure
RestartSec=5
Environment=NODE_DISABLE_COMPILE_CACHE=1
Environment=OPENCLAW_STATE_DIR=/home/wzgrx/.openclaw

[Install]
WantedBy=default.target
